# deploy/videodown-web.service
[Unit]
Description=VideoDown Web Service
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/project/videodown
Environment=PATH=/home/admin/.deno/bin:/usr/local/bin:/usr/bin:/bin
Environment=VIDEODOWN_CONFIG_DIR=/home/project/videodown/config
ExecStart=/home/project/videodown/.venv/bin/uvicorn videodown.web.server:app --host 127.0.0.1 --port 8765 --timeout-keep-alive 300
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
